# autostart/camera.desktop
[Desktop Entry]
Type=Application
Name=camera
Comment=Python Program
Exec=sudo -H -u pi /home/pi/miniconda3/bin/python /home/pi/rpi_img.py &